# libs/gamepad_config/xgo-gamepad.service
[Unit]
Description=XGO Gamepad Controller (evdev → xgolib)
Documentation=man:systemd.service(5)
# 等蓝牙就绪（Xbox 手柄需要蓝牙连接）
After=bluetooth.target xgo-gamepad-config.service
Wants=bluetooth.target
# 注意：若 AI 语音系统也在运行（xgo-rider.service），
#       两者都会使用 xgolib 串口，可能冲突。
#       建议二选一启用，或在主程序中集成本控制器逻辑。

[Service]
Type=simple
User=luwu
# input 组可读取 /dev/input/event*（evdev）
# dialout 组可访问 /dev/ttyAMA* 串口（xgolib）
SupplementaryGroups=input dialout
WorkingDirectory=/home/luwu/XGO-Rider
ExecStart=/usr/bin/python3 /home/luwu/XGO-Rider/gamepad_controller.py
# 手柄未连接时无害重启
Restart=on-failure
RestartSec=5
# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xgo-gamepad

[Install]
WantedBy=multi-user.target
